# scripts/crontab
# Roanes Kitchen Backup Schedule
# Runs daily at 2:00 AM server time
0 2 * * * /backup.sh >> /backups/backup.log 2>&1
